# unsplash/unsplash.service
[Unit]
Description=Apply random unsplash wallpaper

[Service]
Type=simple
ExecStart=/bin/unsplash.sh
